# crontab
0 1 * * * . /etc/environment; python /opt/qbit-maid.py >> /logfile